# supervisord_racebot.conf
[program:racebot]
command = /home/pi/racebot
directory = /tmp
user = root
autostart = true
autorestart = true
stdout_logfile = /var/log/racebot.log
stderr_logfile = /var/log/racebot_err.log

